# chat-page-backup-20260427-085803/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for AskDusk
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel run AskDusk
Restart=on-failure
RestartSec=5s
User=root

[Install]
WantedBy=multi-user.target
